# autostart/my-autostart.desktop
[Desktop Entry]
Name=my auto start
Exec=/home/sam/.config/autostart/autostart.sh
Type=Application
X-GNOME-Autostart-enabled=true
